# upstart/diaspora-sidekiq.conf
start on starting diaspora
stop on stopping diaspora
respawn

exec su -l diaspora -c 'cd /opt/diaspora; RAILS_ENV=production bundle exec sidekiq >> /var/log/diaspora/sidekiq.log 2>&1'
